# rootfs/etc/s6-overlay/s6-rc.d/crond/run
#!/usr/bin/execlineb -P

# run crontab on foreground
crond -f
