# config/crontab
1 0 * * * /backup-r2.sh
3 0 * * * /backup-github.sh
5 0 * * * /backup-rclone.sh