# config/crontab
5 0 * * * /backup-r2.sh
0 0 * * * /backup-github.sh